# docker/supervisord.conf
; Runs the two inference-server processes side by side. Each writes its own
; logfile under /app/logs, which is bind-mounted to ./logs on the host so the
; log-tailer service (and `make docker-logs`) can follow them.
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Observer / metrics agent (FastAPI on :9100) ──────────────────────────────
[program:agent]
command=python -m src.observer_agent.server
directory=/app
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/app/logs/agent.log
stdout_logfile_maxbytes=0
redirect_stderr=true

; ── Inference server (vLLM, or the HF baseline) on :8000 ─────────────────────
; `src.cli serve` exec()s into vLLM, so supervisord supervises vLLM directly.
; SCENARIO is injected by docker-compose (defaults to baseline in the image).
[program:vllm]
command=python -m src.cli serve --scenario %(ENV_SCENARIO)s
directory=/app
autostart=true
autorestart=true
startsecs=10
stopasgroup=true
killasgroup=true
stopwaitsecs=30
stdout_logfile=/app/logs/vllm.log
stdout_logfile_maxbytes=0
redirect_stderr=true
